# cookbooks/aws-parallelcluster-environment/templates/ec2_udev_rules/ec2blkdev.service.erb
[Unit]
Description=Setup ec2blkdev service
After=<%= @network_target %>

[Service]
ExecStart=/sbin/udevadm trigger --action=change --subsystem-match=block; sleep 15

[Install]
WantedBy=multi-user.target
